# license-server/deploy/hackswipe-license.service
[Unit]
Description=HackSwipe License Server
After=network.target

[Service]
Type=simple
# Adjust User for your server (root or a dedicated deploy user).
User=root
WorkingDirectory=/var/www/hackswipe/license-server
EnvironmentFile=/var/www/hackswipe/license-server/.env
# Use the full path to node on the server (which node).
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# Health check helper (optional external monitor can curl this endpoint)
# curl -fsS http://127.0.0.1:3847/health

[Install]
WantedBy=multi-user.target
