# scripts/hackswipe-deploy-webhook.service
[Unit]
Description=HackSwipe GitHub Deploy Webhook
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/hackswipe
EnvironmentFile=-/var/www/hackswipe/scripts/deploy-webhook.env
ExecStart=/root/.nvm/versions/node/v22.14.0/bin/node /var/www/hackswipe/scripts/github-deploy-webhook.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
